# setup 24.04/04-http/etc/cron.d/hew-admin-web
# /etc/cron.d/hew-admin-web: crontab pour la gestion des stats web + sauvegarde mysql
MAILTO=tech@homeasyweb.net

59 23 * * * root /home/hew/bin/GestionLogHttp.php
