# setup 24.04/07-fetchmail/etc/cron.d/hew-fetchmail
# /etc/cron.d/hew-init-fetchmail: crontab pour l'execution des fetchmail

* * * * * hew /home/fetchmail/ex_fetchmail >/dev/null 2>&1
